# mate-display-manager/data/greeter-autostart/gok.desktop.in
[Desktop Entry]
_Name=MATE On-Screen Keyboard
_Comment=Use an on-screen keyboard
TryExec=gok
Exec=gok --login --access-method=directselection
Terminal=false
Type=Application
StartupNotify=true
Categories=MATE;GTK;Accessibility;
AutostartCondition=MATE /desktop/mate/applications/at/screen_keyboard_enabled
X-MATE-Provides=screen-keyboard
X-MATE-Bugzilla-Bugzilla=MATE
X-MATE-Bugzilla-Product=gok
X-MATE-Bugzilla-Component=general

