# mate-display-manager/data/greeter-autostart/mate-mag.desktop.in
[Desktop Entry]
_Name=MATE Screen Magnifier
_Comment=Magnify parts of the screen
TryExec=magnifier
Exec=magnifier -v -m
Terminal=false
Type=Application
StartupNotify=true
Categories=MATE;GTK;Accessibility;
AutostartCondition=MATE /desktop/mate/applications/at/screen_magnifier_enabled
X-MATE-Provides=screen-magnifier
X-MATE-Bugzilla-Bugzilla=MATE
X-MATE-Bugzilla-Product=mate-mag
X-MATE-Bugzilla-Component=general

