# mate-display-manager/data/greeter-autostart/orca-screen-reader.desktop.in
[Desktop Entry]
_Name=Orca Screen Reader
_Comment=Present on-screen information as speech or braille
TryExec=orca
Exec=orca --no-setup --disable main-window --disable splash-window --disable magnifier --enable speech
Icon=orca
Terminal=false
Type=Application
StartupNotify=true
Categories=MATE;GTK;Accessibility;
AutostartCondition=MATE /desktop/mate/applications/at/screen_reader_enabled
X-MATE-Provides=screen-reader
X-MATE-Bugzilla-Bugzilla=MATE
X-MATE-Bugzilla-Product=orca
X-MATE-Bugzilla-Component=general

